# 04_zrok/zrok-mine-v6.service
[Unit]
Description=zrok Share - mine-v6 (port 3000 OneAPI)
After=zrok-agent.service
Wants=zrok-agent.service

[Service]
Type=simple
ExecStart=/usr/bin/script -q -c "/usr/local/bin/zrok2 share public --subordinate -b proxy -n public:mine-v6 http://127.0.0.1:3000" /dev/null
Restart=always
RestartSec=15
User=root

[Install]
WantedBy=multi-user.target
